# pd/tomcat/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/tomcat
drwxr-xr-x	root/tomcat	etc/tomcat/
drwxrwxr-x	root/tomcat	etc/tomcat/Catalina/
-rw-r-----	root/tomcat	etc/tomcat/catalina.policy
-rw-r-----	root/tomcat	etc/tomcat/catalina.properties
-rw-r-----	root/tomcat	etc/tomcat/context.xml
-rw-r-----	root/tomcat	etc/tomcat/logging.properties
-rw-r-----	root/tomcat	etc/tomcat/server.xml
-rw-r-----	root/tomcat	etc/tomcat/tomcat-users.xml
-rwxr-xr-x	root/tomcat	etc/tomcat/tomcat.env
-rw-r-----	root/tomcat	etc/tomcat/web.xml
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/jsvc
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/apr-1/
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libtcnative-1.a
-rwxr-xr-x	root/root	usr/lib/libtcnative-1.la
lrwxrwxrwx	root/root	usr/lib/libtcnative-1.so -> libtcnative-1.so.0.1.33
lrwxrwxrwx	root/root	usr/lib/libtcnative-1.so.0 -> libtcnative-1.so.0.1.33
-rwxr-xr-x	root/root	usr/lib/libtcnative-1.so.0.1.33
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/tcnative-1.pc
drwxr-xr-x	root/root	usr/lib/tomcat/
-rw-r--r--	root/root	usr/lib/tomcat/annotations-api.jar
drwxr-xr-x	root/root	usr/lib/tomcat/bin/
-rw-r--r--	root/root	usr/lib/tomcat/bin/bootstrap.jar
-rw-r--r--	root/root	usr/lib/tomcat/bin/catalina-tasks.xml
-rw-r--r--	root/root	usr/lib/tomcat/bin/commons-daemon.jar
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/shutdown.sh
-rwxr-xr-x	root/root	usr/lib/tomcat/bin/startup.sh
-rw-r--r--	root/root	usr/lib/tomcat/bin/tomcat-juli.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina-ant.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina-ha.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina-tribes.jar
-rw-r--r--	root/root	usr/lib/tomcat/catalina.jar
lrwxrwxrwx	root/root	usr/lib/tomcat/conf -> ../../../etc/tomcat
-rw-r--r--	root/root	usr/lib/tomcat/ecj-4.4.2.jar
-rw-r--r--	root/root	usr/lib/tomcat/el-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/jasper-el.jar
-rw-r--r--	root/root	usr/lib/tomcat/jasper.jar
-rw-r--r--	root/root	usr/lib/tomcat/jsp-api.jar
lrwxrwxrwx	root/root	usr/lib/tomcat/lib -> .
lrwxrwxrwx	root/root	usr/lib/tomcat/logs -> ../../../var/log/tomcat
-rw-r--r--	root/root	usr/lib/tomcat/servlet-api.jar
lrwxrwxrwx	root/root	usr/lib/tomcat/temp -> ../../../var/run/tomcat/temp
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-api.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-coyote.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-dbcp.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-jdbc.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat-util.jar
-rw-r--r--	root/root	usr/lib/tomcat/tomcat7-websocket.jar
lrwxrwxrwx	root/root	usr/lib/tomcat/webapps -> ../../../var/tomcat/webapps
-rw-r--r--	root/root	usr/lib/tomcat/websocket-api.jar
lrwxrwxrwx	root/root	usr/lib/tomcat/work -> ../../../var/tomcat/work
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	tomcat/tomcat	var/log/tomcat/
drwxr-xr-x	root/root	var/tomcat/
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/ROOT/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/RELEASE-NOTES.txt
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/ROOT/WEB-INF/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/WEB-INF/web.xml
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/asf-logo-wide.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/asf-logo.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/bg-button.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/bg-middle.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/bg-nav-item.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/bg-nav.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/bg-upper.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/build.xml
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/favicon.ico
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/index.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/tomcat-power.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/tomcat.css
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/tomcat.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/tomcat.png
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/ROOT/tomcat.svg
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/host-manager/
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/host-manager/META-INF/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/META-INF/context.xml
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/host-manager/WEB-INF/
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/host-manager/WEB-INF/jsp/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/WEB-INF/jsp/401.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/WEB-INF/jsp/403.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/WEB-INF/jsp/404.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/WEB-INF/web.xml
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/host-manager/images/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/add.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/asf-logo.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/code.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/design.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/docs.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/fix.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/tomcat.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/update.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/images/void.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/index.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/host-manager/manager.xml
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/manager/
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/manager/META-INF/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/META-INF/context.xml
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/jsp/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/jsp/401.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/jsp/403.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/jsp/404.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/jsp/sessionDetail.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/jsp/sessionsList.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/WEB-INF/web.xml
drwxr-xr-x	tomcat/tomcat	var/tomcat/webapps/manager/images/
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/add.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/asf-logo.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/code.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/design.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/docs.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/fix.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/tomcat.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/update.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/images/void.gif
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/index.jsp
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/status.xsd
-rw-r--r--	tomcat/tomcat	var/tomcat/webapps/manager/xform.xsl
drwxr-xr-x	tomcat/tomcat	var/tomcat/work/
